# anki-watcher.service
[Unit]
Description=Anki Watcher — Syncs vocabulary cards from Syncthing to Anki
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/.config/anki-watcher/anki_watcher.py --config %h/.config/anki-watcher/config.json
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
